# porthos/root/etc/cron.d/porthos.cron
#
# porthos crontab
#

RANDOM=20
CRON_TZ=UTC

0 21 * * 0-4      root /usr/local/bin/repo-bulk-pull
0 22 * * 0-4      root /usr/local/bin/repo-bulk-cleanup >/dev/null

# end

